# rclone/.config/systemd/user/rclone-bisync.timer
[Unit]
Description=Periodic rclone bisync

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Unit=rclone-bisync.service

[Install]
WantedBy=timers.target
